# docker/runit/simple/run
#!/bin/sh
# vim: set ft=sh:

exec 2>&1
exec /usr/bin/ss-server -c /etc/shadowsocks-simple.json
