# ansible/roles/9seconds-docker-shadow/templates/shadow.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=Shadowsocks

[Service]
ExecStart=/bin/bash /usr/share/docker-shadow.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
