# docker/runit/kcptun/run
#!/bin/sh
# vim: set ft=sh:

exec 2>&1
exec /usr/bin/kcp-server -c /etc/shadowsocks-simple.json
